# deploy/arcturus-bus.service
[Unit]
Description=Arcturus message bus
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/arcturus/platform.env
ExecStartPre=/usr/bin/mkdir -p %t/arcturus
ExecStart=/usr/bin/node %h/.local/share/arcturus-deployer/current/modules/bus/dist/index.js
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t/arcturus

[Install]
WantedBy=default.target
