# deploy/arcturus-registry.service
[Unit]
Description=Arcturus manifest registry
After=arcturus-bus.service
Wants=arcturus-bus.service

[Service]
Type=simple
EnvironmentFile=%h/.config/arcturus/platform.env
ExecStartPre=/usr/bin/mkdir -p %t/arcturus
ExecStart=/usr/bin/node %h/.local/share/arcturus-deployer/current/modules/registry/dist/index.js
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t/arcturus %h/.local/share/arcturus-deployer/active-manifests

[Install]
WantedBy=default.target
